# crontab
# sample
* * * * * root echo hoge | logger
